# deploy/systemd/dailywall-crawl-test.service
[Unit]
Description=DailyWall crawl job for 72-hour systemd test
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=ops
Group=ops
WorkingDirectory=/home/ops/Projects/DailyWall
ExecStart=/usr/bin/env bash /home/ops/Projects/DailyWall/scripts/run_crawl_job.sh systemd-test
Nice=10

