# deploy/systemd/dailywall-api.service
[Unit]
Description=DailyWall FastAPI service
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=ops
Group=ops
WorkingDirectory=/home/ops/Projects/DailyWall
Environment=PYTHONUNBUFFERED=1
ExecStartPre=/usr/bin/test -x /home/ops/Projects/DailyWall/.venv/bin/python
ExecStart=/home/ops/Projects/DailyWall/.venv/bin/python -m app.main
Restart=on-failure
RestartSec=5
SyslogIdentifier=dailywall-api

[Install]
WantedBy=multi-user.target
